# static/tests/data/supervisord/simple/supervisord.conf
[program:mock-awesome]
command = java -cp "/var/wiremock/lib/*:/var/wiremock/extensions/*" com.github.tomakehurst.wiremock.standalone.WireMockServerRunner --port 8000 --root-dir static/tests/data/wiremock/configs/two-services/awesome
