# static/tests/data/supervisord/two-services/supervisord.conf
[include]
files = mocks/*.conf

[program:mocks-watcher]
environment = A="env 1",B="this is a test"
command = watcher --mocks="hello"
