# static/tests/data/supervisord/one-service/mocks/mock-awesome.conf
[program:mock-awesome]
environment = ROOT="/home/mock/awesome",PORT="8000",NAME="awesome"
command = java -cp "/var/wiremock/lib/*:/var/wiremock/extensions/*" com.github.tomakehurst.wiremock.standalone.WireMockServerRunner --port 8000 --root-dir /static/tests/data/wiremock/configs/two-services/awesome
